# packaging/holoscan-proxy.service
[Unit]
Description=Holoscan Proxy Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/holoscan-proxy --config /etc/holoscan-proxy/config.yaml
Restart=on-failure
RestartSec=2
User=nobody
WorkingDirectory=/var/lib/holoscan-proxy
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

